# packer_files/webapp.service
[Unit]
Description=Web Application Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/tmp/webapp
ExecStart=/usr/bin/npm start
Restart=on-failure

[Install]
WantedBy=multi-user.target